# linux-systemd.service.example
[Unit]
Description=RMate Launcher - Remote file editing server
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rmate_launcher
Restart=always
RestartSec=5

# Environment variables
# Be sure to specify the full path to the editor in RMATE_EDITOR since the regular $PATH is not set
Environment=RMATE_EDITOR=/usr/local/bin/code --wait
# RMATE_SOCKET defaults to %h/.rmate_launcher/rmate.sock if not specified
# Environment=RMATE_SOCKET=%h/.rmate_launcher/rmate.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rmate_launcher

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=tmpfs
ReadWritePaths=%h/.rmate_launcher
ReadOnlyPaths=%h/.config %h/.local
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target