# linux-systemd.service.example
[Unit]
Description=RMate Server - Remote file editing server
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rmate_server
Restart=always
RestartSec=5

# Environment variables
Environment=RMATE_EDITOR=code --wait
# RMATE_SOCKET defaults to %h/.rmate-server/rmate.sock if not specified
# Environment=RMATE_SOCKET=%h/.rmate-server/rmate.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rmate-server

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=tmpfs
ReadWritePaths=%h/.rmate-server
ReadOnlyPaths=%h/.config %h/.local
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target