# linux-systemd.service.example
[Unit]
Description=RMate Launcher - Remote file editing server
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rmate_launcher
Restart=on-failure
RestartSec=5

# Environment variables
# Be sure to specify the full path to the editor in RMATE_EDITOR since the regular $PATH is not set
Environment=RMATE_EDITOR=/usr/local/bin/code --wait
# RMATE_SOCKET defaults to %h/.rmate_launcher/rmate.sock if not specified
# Environment=RMATE_SOCKET=%h/.rmate_launcher/rmate.sock

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rmate_launcher

[Install]
WantedBy=graphical-session.target