# macos-launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.rmate-server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/rmate_server</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>RMATE_EDITOR</key>
        <string>code --wait</string>
        <!-- RMATE_SOCKET defaults to ~/.rmate-server/rmate.sock if not specified -->
        <!-- 
        <key>RMATE_SOCKET</key>
        <string>/Users/USERNAME/.rmate-server/rmate.sock</string>
        -->
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/rmate-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/rmate-server-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/.rmate-server</string>
    
    <!-- Note: Update USERNAME above to your actual username -->
    
    <!-- Uncomment the sandbox profile below for additional security -->
    <!-- 
    <key>SandboxProfile</key>
    <string>
        (version 1)
        (deny default)
        (allow process-exec (literal "/usr/local/bin/rmate_server"))
        (allow file-read* file-write* (subpath "/Users/USERNAME/.rmate-server"))
        (allow file-read* (subpath "/usr/bin") (subpath "/usr/local/bin"))
        (allow file-read* (subpath "/Applications"))
        (allow network-outbound (remote ip))
        (allow network-bind (local ip))
    </string>
    -->
</dict>
</plist>